# php-fpm/7.2-zts/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm

[program:cron]
command=/usr/sbin/cron -f
